# deploy/weixin-session-manager.service
[Unit]
Description=weixin-agent session manager
After=weixin-agent.service
Requires=weixin-agent.service
StartLimitIntervalSec=300
StartLimitBurst=6

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/weixin-agent
ExecStart=/home/ubuntu/weixin-agent/venv/bin/python /home/ubuntu/weixin-agent/session-manager.py
Restart=always
RestartSec=10
User=ubuntu
EnvironmentFile=/home/ubuntu/weixin-agent/.env
Environment=TZ=Asia/Shanghai
NoNewPrivileges=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=weixin-session-mgr

[Install]
WantedBy=multi-user.target
